# deploy/rapidboxes.service
[Unit]
Description=RaPiD-boxes plant imaging chamber controller
After=network.target

[Service]
Type=simple
User=@USER@
WorkingDirectory=@BACK_DIR@
EnvironmentFile=/etc/rapidboxes.env
ExecStart=@VENV@/bin/python -m rapidboxes
Restart=always
RestartSec=1
# Hard-kill promptly: graceful stop can hang on open WebSockets / MJPEG.
TimeoutStopSec=5
KillMode=control-group
KillSignal=SIGKILL

[Install]
WantedBy=multi-user.target
